# systemd.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/eth-log-alerting/service.sh

[Install]
WantedBy=default.target